# ora-initial-setup/DEBIAN/prerm
#!/bin/bash

if [ "$1" == "remove" ]; then
    /usr/sbin/update-rc.d oracle-shm remove
fi
